# writ/package/install.sh
#!/bin/bash
chmod 755 /etc/init.d/writ
update-rc.d writ defaults